# fix-files/package/base-files/files/etc/uci-defaults/12_network-generate-ula
grep poweroff /etc/config/luci >/dev/null || cat <<\EOF >>/etc/config/luci

config command
	option name '关闭系统'
	option command 'poweroff'

config command
	option name '网络接口'
	option command 'ifconfig'
EOF

[ "$(uci -q get network.globals.ula_prefix)" != "auto" ] && exit 0

uci -q batch <<\EOF >/dev/null
	set network.globals.ula_prefix=fd00::/48
	commit network
	set upnpd.config.enabled=1
	set upnpd.config.uuid=3943084a-dd61-4636-989c-ab0123456789
	commit upnpd
EOF

exit 0
